# system/DC-detection-gpio.service
[Unit]
Description=DC-detection-gpio
Requires=network.target
 
[Service]
Type=simple
User=root
ExecStart=/bin/bash /var/www/html/DC-detection-gpio.sh
Restart=always
RestartSec=5
 
[Install]
WantedBy=multi-user.target
